# systemd/ai-helper.service
[Unit]
Description=AI Helper — Personal On-Device Automation Tool
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/ai-helper-app
ExecStart=/bin/bash -c 'export NVM_DIR="${HOME}/.config/nvm"; [ -s "${NVM_DIR}/nvm.sh" ] || NVM_DIR="${HOME}/.nvm"; source "${NVM_DIR}/nvm.sh"; cd %h/ai-helper-app && pnpm start'
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production
Environment=AI_HELPER_LIVE=1

[Install]
WantedBy=default.target
